# kiosk.service
[Unit]
Description=Kiosk Display Service
After=network.target graphical.target

[Service]
Type=simple
User=match
WorkingDirectory=/home/match/studitv-kiosk
ExecStart=/home/match/studitv-kiosk/launch_kiosk.sh
Restart=always
RestartSec=5
WatchdogSec=30
Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target